# src/supervisord.conf
[supervisord]
nodaemon=true
[program:xbmc-server]
command=/usr/local/lib/xbmc/xbmc-server.bin --no-test --nolirc -p
